# roles/pangolin-newt/templates/newt.service.j2
[Unit]
Description=Newt VPN Client
After=network.target

[Service]
# Zugangsdaten liegen in einer nur für root lesbaren EnvironmentFile statt auf
# der ExecStart-Zeile: dort wären sie für alle lokalen Nutzer sichtbar (Unit-
# Datei ist 0644, Kommandozeile via ps einsehbar).
EnvironmentFile=/etc/newt/newt.env
ExecStart=/usr/local/bin/newt
Restart=always
User=root

[Install]
WantedBy=multi-user.target
